# event_loop_histogram_c/ngx-bpf-stats/etc/systemd/system/ngx-bpf-stats.service
[Unit]
Description=Nginx eBPF statistics exporter

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ngx-bpf-stats
